# scripts/qwenpaw_bridge_supervisor.conf
[program:qwenpaw-miloco-bridge]
command=python3 /path/to/qwenpaw_webhook_bridge.py
autorestart=true
startretries=3
startsecs=2
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=%(ENV_HOME)s/.openclaw/miloco/log/qwenpaw_bridge.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=QWENPAW_BIN="/app/venv/bin/qwenpaw",QWENPAW_BASE_URL="http://localhost:8088"
